# crontab
1 1 * * *   find $HOME/scratch/ -exec touch -a {} \;
2 2 * * *   rsync -avz --delete /scratch/ng1154e/ $HOME/fsnreg/scratch/
